# roles/bootstrap/templates/supervisord.service.j2
[Unit]
Description=Supervisord Init script

[Service]
WorkingDirectory=/root
Type=forking
ExecStart=/usr/bin/supervisord
KillMode=process

[Install]
WantedBy=multi-user.target
